# dist/src/main/resources/launchd/scripts/postinstall
#!/bin/sh

set -e
/bin/launchctl load /Library/LaunchDaemons/org.ogerardin.back2back.plist
